# config/roles/upgrade/files/etc/systemd/system/download-updates.service
[Unit]
Description=Download available updates

[Service]
Type=oneshot
ExecStart=/usr/local/bin/download-updates

[Install]
WantedBy=basic.target
