# config/roles/upgrade/files/etc/systemd/system/download-updates.timer
[Unit]
Description=Download updates daily

[Timer]
OnCalendar=0:00:00

[Install]
WantedBy=timers.target
