# config/roles/upgrade/files/etc/systemd/system/pacman-keyring-update.service
[Unit]
Description=Keep pacman keyring updated

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman-key --refresh-keys

[Install]
WantedBy=basic.target
