# config/roles/upgrade/files/etc/systemd/system/pacman-keyring-update.timer
[Unit]
Description=Update pacman keyring once a day

[Timer]
OnCalendar=0:00:00

[Install]
WantedBy=timers.target
