# config/roles/upgrade/files/etc/systemd/system/safe-upgrade.service
[Unit]
Description=Install updates on-demand

[Service]
Type=oneshot
ExecStart=/usr/local/bin/safe-upgrade

[Install]
WantedBy=basic.target
