# config/roles/upgrade/files/etc/systemd/system/safe-upgrade.timer
[Unit]
Description=Attempt installing updates

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
