# config/roles/upgrade/files/etc/systemd/system/schedule-reboot.path
[Unit]
Description=Reboot for updates on demand

[Path]
PathExists=/system-update

[Install]
WantedBy=basic.target
